# etc/profile.d/custom_environment_variables.sh
export EDITOR=/usr/bin/vim
export VISUAL=/usr/bin/gedit
export DIFFPROG=/usr/bin/meld
export BROWSER=/usr/bin/google-chrome-stable
export QT_QPA_PLATFORMTHEME=qt5ct
